# PortfolioBackend/scripts/portfolio-backend.service
[Unit]
Description=Portfolio Backend ASP.NET Core App
After=network.target
Wants=network-online.target

[Service]
Type=notify
# Change to your app directory
WorkingDirectory=/opt/portfolio-backend
ExecStart=/opt/portfolio-backend/PortfolioBackend
Restart=always
RestartSec=30
KillSignal=SIGINT
SyslogIdentifier=portfolio-backend
User=www-data
Group=www-data

# Environment variables
Environment=ASPNETCORE_ENVIRONMENT=Production
Environment=ASPNETCORE_URLS=https://+:443;http://+:80
Environment=DOTNET_PRINT_TELEMETRY_MESSAGE=false

# SSL Certificate environment variables (set these via systemctl edit)
Environment=SSL_CERT_PASSWORD=
Environment=EMAIL_PASSWORD=

# Security settings
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/opt/portfolio-backend/data
ReadWritePaths=/opt/portfolio-backend/ssl
ReadWritePaths=/opt/portfolio-backend/uploads
ReadWritePaths=/opt/portfolio-backend/logs

# Process limits
LimitNOFILE=1048576
LimitNPROC=1048576

[Install]
WantedBy=multi-user.target
