# deployment/codenex.service
[Unit]
Description=CodeNex .NET Web Application
After=network.target

[Service]
Type=notify
# Run as a specific user (recommended for security)
User=www-data
Group=www-data

# Application directory - UPDATE THIS PATH
WorkingDirectory=/var/www/codenex
ExecStart=/usr/bin/dotnet /var/www/codenex/CodeNex.dll

# Restart configuration
Restart=always
RestartSec=10
KillSignal=SIGINT
SyslogIdentifier=codenex
Environment=ASPNETCORE_ENVIRONMENT=Production
Environment=DOTNET_PRINT_TELEMETRY_MESSAGE=false

# Performance and security settings
PrivateTmp=true
NoNewPrivileges=true

# Logging
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
